# scripts/99-reload-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/bin/nohup /home/yottaphy/.config/scripts/monitorchange.sh > /dev/null 2>&1 &"
